# .config/systemd/user/motd.service
[Service]
Type=oneshot
ExecStart=/home/kay/Scripts/motd.sh
