# .config/systemd/user/motd.timer
[Unit]
Description=Refresh the motd periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
